# stage2-zf/01-sys-tweaks/files/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin zf --noclear %I $TERM